# docker/supervisor.d/cassandra.ini
[program:cassandra]
user=reactive-sandbox
command=/opt/cassandra/bin/cassandra -f
autostart=true
autorestart=true

