# docker/supervisor.d/elasticsearch.ini
[program:elasticsearch]
user=reactive-sandbox
command=/opt/elasticsearch/bin/elasticsearch
autostart=true
autorestart=true

