# docker/supervisor.d/kafka.ini
[program:kafka]
user=reactive-sandbox
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
autostart=true
autorestart=true

